# src/etc/systemd/system/prepare_livesystem.service
[Unit]
Description=Prepare the BBQLinux livesystem

[Service]
ExecStart=/usr/bin/prepare_livesystem

[Install]
WantedBy=multi-user.target 
